# systemd/user/mbsync@.timer
[Unit]
Description=mbsync for %i

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
